# prometheus_process_exporter/rootfs/etc/services.d/process_exporter/run
#!/usr/bin/env bashio
# ==============================================================================
# Home Assistant Community Add-on: Prometheus process Exporter
# Runs the Prometheus process Exporter
# ==============================================================================
bashio::log.info "Starting Prometheus Process Exporter..."

# Run Prometheus
exec /usr/local/bin/process_exporter --web.config.file=/config/prometheus_process_exporter/process_exporter_web.yml
