# prometheus_process_exporter/rootfs/etc/services.d/process_exporter/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Prometheus Process Exporter
# Runs the Prometheus Process Exporter
# ==============================================================================
bashio::log.info "Starting prometheus process exporter..."

exec ls -alh /usr/local/bin/process-exporter

# Run Prometheus
exec s6-setuidgid prometheus /usr/local/bin/process-exporter -config.path=/etc/prometheus_process-exporter.yml

# Debug
#exec sleep 3000
